# bot/systemd/vadana-bot.service
[Unit]
Description=Vadana slides Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vadana-extractor
EnvironmentFile=/opt/vadana-extractor/bot/.env
ExecStart=/opt/vadana-extractor/venv/bin/python -m bot.bot
Restart=always
RestartSec=5
MemoryMax=1G
CPUQuota=150%

[Install]
WantedBy=multi-user.target
